# root/etc/services.d/shiori/run
#!/usr/bin/with-contenv bash

export SHIORI_DIR="/config"

exec s6-setuidgid abc /usr/local/bin/shiori serve
